# templates/default/flinkhistoryserver.service.erb
[Unit]
Description = Flink historyserver
After = syslog.target network.target remote-fs.target <%= @deps %>
#restart every 5s and at most 121 times in 11min to guaranty that it will restart for at least 10min (120*5s)
StartLimitIntervalSec=660
StartLimitBurst=121

[Service]
User = <%= node['flink']['user'] %>
EnvironmentFile=<%= node['flink']['historyserver']['environment'] %>
PIDFile = /tmp/flink-<%= node['flink']['user'] %>-historyserver.pid
ExecStart = <%= node['flink']['base_dir'] %>/bin/historyserver.sh start
ExecStop = <%= node['flink']['base_dir'] %>/bin/historyserver.sh stop
LimitNOFILE = 50000
Type = forking
Restart=on-failure
RestartSec=5

[Install]
WantedBy = multi-user.target
