# templates/default/taskmanager.service.erb
[Unit]
Description = Taskmanager, worker for Flink.
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.flink.user %>

PIDFile = <%= node.flink.logs_dir %>/hadoop-<%= node.apache_hadoop.hdfs.user %>-namenode.pid
ExecStart = <%= node.flink.home %>/sbin/start-nn.sh
ExecStop = <%= node.flink.home %>/sbin/stop-nn.sh

Type = forking

[Install]
WantedBy = multi-user.target
